# rc.d/tun2socks
#!/bin/sh

# PROVIDE: tun2socks
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="tun2socks"
rcvar="tun2socks_enable"
pidfile="/var/run/${name}.pid"
logfile="/var/log/${name}.log"
command="/usr/local/bin/tun2socks"
command_args="--config /usr/local/etc/tun2socks/config.yaml"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
restart_cmd="${name}_restart"
status_cmd="${name}_status"


tun2socks_start() {
    if [ -f "$pidfile" ] && kill -0 $(cat "$pidfile") 2>/dev/null; then
        echo "${name} is already running." | tee -a "$logfile"
        return 1
    fi
    echo "$(date '+%Y-%m-%d %H:%M:%S') - Starting ${name}..." | tee -a "$logfile"
    nohup ${command} ${command_args} >> "$logfile" 2>&1 &
    echo $! > "$pidfile"
    echo "$(date '+%Y-%m-%d %H:%M:%S') - ${name} started successfully." | tee -a "$logfile"
}

tun2socks_stop() {
    if [ ! -f "$pidfile" ] || ! kill -0 $(cat "$pidfile") 2>/dev/null; then
        echo "${name} is not running." | tee -a "$logfile"
        return 1
    fi
    echo "$(date '+%Y-%m-%d %H:%M:%S') - Stopping ${name}..." | tee -a "$logfile"
    kill $(cat "$pidfile") && rm -f "$pidfile"
    echo "$(date '+%Y-%m-%d %H:%M:%S') - ${name} has stopped." | tee -a "$logfile"
}

tun2socks_restart() {
    echo "$(date '+%Y-%m-%d %H:%M:%S') - Restarting ${name}..." | tee -a "$logfile"
    tun2socks_stop
    sleep 1
    tun2socks_start
}

tun2socks_status() {
    if [ -f "$pidfile" ] && kill -0 $(cat "$pidfile") 2>/dev/null; then
        echo "tun2socks is running (PID: $(cat "$pidfile"))."
    else
        echo "tun2socks is not running."
    fi
}
load_rc_config $name
run_rc_command "$1"